# supervisord/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:squid]
# Start the squid proxy
# -N: Run in the foreground (critical for supervisor)
# -f: Point to our custom config file
command=/usr/sbin/squid -N -f /etc/squid/squid.conf
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tc-ui]
# Start our Go web application
command=/usr/local/bin/tc-ui
# Run it from the /app directory where the UI files are
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
